# systemd/tgbot.service
# Example configuration for running the bot on a server as a systemd service
[Unit]
Description=OrendaUAbot daemon
After=network.target

[Service]
User=tgbot
Group=tgbot
Type=simple
WorkingDirectory=/opt/OrendaUAbot
ExecStart=/opt/OrendaUAbot/venv/bin/python /opt/OrendaUAbot/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
