# systemd/tgbot.service
# Приклад конфігурації для запуску бота на сервері в якості системної служби
[Unit]
Description=OrendaUAbot daemon
After=network.target

[Service]
User=tgbot
Group=tgbot
Type=simple
WorkingDirectory=/opt/OrendaUAbot
ExecStart=/opt/OrendaUAbot/venv/bin/python /opt/OrendaUAbot/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
